# share/scripts/crontab
*/2 * * * * /home/john/tpwater/tpwater.sh   start
45   2 * * * /home/john/tpwater/tpwater.sh restart
